# scripts/com.borjatalavera.wcpredictor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.borjatalavera.wcpredictor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/borjatalavera/Desktop/AlgoritmoPredict/scripts/update_demo.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/borjatalavera/Library/Logs/wcpredictor.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/borjatalavera/Library/Logs/wcpredictor.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
